# backend/deploy/glucose-api.service.example
[Unit]
Description=家有糖人 API
After=network.target

[Service]
Type=simple
User=glucose
Group=glucose
WorkingDirectory=/opt/app/glucose-api
EnvironmentFile=/opt/app/glucose-api/.env
ExecStart=/opt/app/glucose-api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8082 --workers 1 --no-access-log
StandardOutput=append:/var/log/glucose-api/app.log
StandardError=append:/var/log/glucose-api/error.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
